# templates/etc/systemd/system/mariadb.service.d/timeout-start-sec.conf.j2

[Service]
# Override defualt timeout for starting service, due to long SST duration
TimeoutStartSec={{ mariadb_timeout_start_sec }}

